# templates/aqi_crawl.service.j2
[Unit]
Description=crawl aqi data

[Service]
ExecStart={{ user_systemd_path }}.local/bin/scrapy crawl aqi
WorkingDirectory={{ user_systemd_path }}/workspace

[Install]
WantedBy=default.target
