# templates/aqi_crawl.service.j2
[Unit]
Description=crawl aqi data

[Service]
ExecStart={{ ansible_env.HOME }}/.local/bin/scrapy crawl aqi
WorkingDirectory={{ proj_path }}

[Install]
WantedBy=default.target
